# systemd_files/mam-seedbox.timer
[Unit]
Description=Update MAM seedbox IP daily

[Timer]
Unit=mam-seedbox.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target